# .opencode/dream.cron
# HeartFlow v8.7.1 定时任务配置
# 项目目录：/Users/apple/mark-heartflow-claw

# 每 24 分钟增强理论升级
*/24 * * * * /Users/apple/mark-heartflow-claw/scripts/hourly-enhanced-upgrade-v2.sh >> /Users/apple/mark-heartflow-claw/logs/enhanced-cron.log 2>&1
